# services/bearbox.service
[Unit]
Description=BearBox
After=multi-user.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /home/bearbox/bearbox/core/profile_manager.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target